# rootfs/etc/services.d/phpfpm/run
#!/usr/bin/execlineb -P

with-contenv
php-fpm7 -F
