# rootfs/etc/services.d/phpfpm/run
#!/command/with-contenv sh

php-fpm84 -F
