# rootfs/etc/services.d/phpfpm/run
#!/usr/bin/with-contenv sh

php-fpm8 -F
