# rootfs/etc/services.d/phpfpm/run
#!/command/with-contenv sh

php-fpm83 -F
