# rootfs/etc/services.d/nginx/run
#!/usr/bin/execlineb -P

with-contenv
nginx -c /etc/nginx/nginx.conf
