# rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv sh

nginx -c /etc/nginx/nginx.conf
