# conf/spotweb.cron
15 * * * *  /usr/bin/php7 /var/www/spotweb/retrieve.php
