# conf/spotweb.cron
15 * * * *  /usr/bin/php /var/www/spotweb/retrieve.php
