# rootfs/etc/periodic/15min/spotweb
#!/bin/sh

s6-setuidgid abc php /app/retrieve.php
